# packaging/macos/tech.fixedwidth.glass.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  gui/<uid> LaunchAgent template that runs the signed GlassMcp.app as
  `glass-mcp serve --http --menubar` — the visible menu-bar app, showing a
  status item the user can Quit from (see docs/how-to/setup-macos.md and
  packaging/macos/README.md).

  KeepAlive is false so a user's Quit sticks: launchd won't relaunch the job
  until the next login (RunAtLoad starts it then). Onboarding's one-shot
  restart (`launchctl kickstart -k`) is how a fresh process gets spawned in
  between, independent of KeepAlive.

  Before loading this, edit the two placeholders below:
    - /Users/YOU  -> your home directory (also fix the app path if you didn't
                     install GlassMcp.app to /Applications).
  The bind address (127.0.0.1:7300) is loopback-only, so no bearer token is
  required; see docs/how-to/run-over-the-network.md for reaching it from another machine.
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>tech.fixedwidth.glass</string>

	<key>ProgramArguments</key>
	<array>
		<string>/Applications/GlassMcp.app/Contents/MacOS/glass-mcp</string>
		<string>serve</string>
		<string>--http</string>
		<string>--menubar</string>
		<string>--addr</string>
		<string>127.0.0.1:7300</string>
	</array>

	<!--
	  A writable working directory. launchd otherwise starts the job with cwd `/` (a
	  read-only sealed volume), which makes glass_baseline_save / glass_diff fail with a
	  read-only-filesystem error. The home dir is always writable.
	-->
	<key>WorkingDirectory</key>
	<string>/Users/YOU</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>GLASS_BACKEND</key>
		<string>macos</string>
	</dict>

	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<false/>
	<key>ProcessType</key>
	<string>Background</string>

	<key>StandardOutPath</key>
	<string>/Users/YOU/Library/Logs/GlassMcp/stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/YOU/Library/Logs/GlassMcp/stderr.log</string>
</dict>
</plist>
